# conf/apache.conf
[program:apache]
command=/usr/sbin/apachectl start
environment=APACHE_ARGUMENTS=-DFOREGROUND
autostart=true
autorestart=true
redirect_stderr=true
