# conf/django_app.template.conf
[program:myDjangoApp]
directory=/path/to/django/project
command=/usr/bin/gunicorn --workers=4 project.wsgi:application
user=john
umask=022
autostart=True
autorestart=True
redirect_stderr=True
