# conf/java_app.template.conf
[program:NAME]
command=java -Xms128m -Xmx1g -XX:+UseCompressedOops -XX:+UseParNewGC -XX:+UseConcMarkSweepGC -jar /opt/NAME.jar
autostart=true
autorestart=true
stderr_logfile=/var/log/NAME.err.log
