# conf/redis.conf
[program:redis]
command=/usr/local/bin/redis-server /etc/redis.conf
autorestart=true
stdout_logfile=/var/log/redis.log
