# conf/znc.conf
[program:znc]
command=/usr/bin/znc -f --datadir=/home/znc/.znc
user=znc
